# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=service mysql start
autorestart=true
