# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=service mysql start
autorestart=true

[program:mariadb_setting]
command=/tmp/allow-remote-connection.sh
autorestart=false
